# packages/fc-env/ubuntu/devbookd.service
[Unit]
Description=Devbook Daemon Service

[Service]
Type=simple
Restart=always
User=root
Group=root
ExecStart=/usr/bin/bash -l -c "/usr/bin/devbookd"

[Install]
WantedBy=multi-user.target
