# packages/fc-envs/alpine/devbookd-init
#!/sbin/openrc-run
name=${RC_SVCNAME}
supervisor="supervise-daemon"
command="/usr/bin/devbookd"
pidfile="/run/devbookd.pid"
#output_log="/var/log/devbookd.log"
#error_log="/var/log/devbookd.err"

depend() {
  need net
  need localmount
  #use logger
}

reload() {
  ebegin "Reloading devbook-daemon"
  ${supervisor} ${RC_SVCNAME} --signal HUP --pidfile "${pidfile}"
  eend $?
}
